# piremd.service
[Unit]
Description=infrared controller for linux

[Service]
Type=simple
ExecStart=/usr/sbin
