# daemon_install_files/piremd.service
[Unit]
Description=Infrared controller for linux. config file: (/etc/piremd/config.json)

[Service]
Type=simple
ExecStart=/usr/sbin/piremd
